# dot_config/systemd/user/backup_restic_local.service
[Unit]
Description=Daily Backup With Restic
 
[Service]
Type=oneshot
ExecStart=/home/tlsh/.restic-backup/backup_local.sh
