# dot_config/systemd/user/backup_cdda.service
[Unit]
Description=Cataclysm DDA backup
 
[Service]
Type=oneshot
ExecStart=/home/tlsh/.cdda_backup/backup_cdda.py
