# dot_config/systemd/user/backup_restic_local.timer
[Unit]
Description=Daily backup Timer
 
[Timer]
OnBootSec=1min
OnUnitActiveSec=1d
 
[Install]
WantedBy=timers.target
