# dot_config/systemd/user/backup_borg.service
[Unit]
Description=Daily Backup
 
[Service]
Type=oneshot
ExecStart=/home/tlsh/.backup/backup_borg.py
