# dot_config/systemd/user/copyq_tab.service
[Unit]
Description=copyq tab tmpfs

[Service]
Type=oneshot
ExecStart=/home/tlsh/.local/bin/copyq_tab.sh

[Install]
WantedBy=default.target
